# scripts/com.crypto-balanced-strategy.auto.adaptive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.crypto-balanced-strategy.auto</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/rongts/.agents/skills/crypto-balanced-strategy/scripts/auto_dual_cycle.py</string>
    <string>--live</string>
    <string>--promote-days</string>
    <string>2</string>
    <string>--allow-aggressive</string>
    <string>--aggressive-promote-days</string>
    <string>5</string>
    <string>--network-recover-retry</string>
    <string>--network-recover-max-wait-minutes</string>
    <string>360</string>
    <string>--aggressive-ratio</string>
    <string>0.3</string>
    <string>--kill-switch-file</string>
    <string>/Users/rongts/.agents/skills/crypto-balanced-strategy/results/kill_switch</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/rongts/.agents/skills/crypto-balanced-strategy</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OKX_API_KEY</key>
    <string>REPLACE_ME</string>
    <key>OKX_API_SECRET</key>
    <string>REPLACE_ME</string>
    <key>OKX_API_PASSPHRASE</key>
    <string>REPLACE_ME</string>
    <key>AUTO_TRADING_ENABLED</key>
    <string>true</string>
  </dict>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/crypto-balanced-strategy-auto.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/crypto-balanced-strategy-auto.err.log</string>
</dict>
</plist>
